# user_data/strategies/smarttsa.service
[Unit]
Description=SmartTSA Bot
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

WorkingDirectory=/home/ubuntu/SmartTSA
ExecStart=/usr/bin/python3 -m freqtrade trade --config /home/ubuntu/SmartTSA/user_data/config.json --strategy SmartTSA

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
